# packaging/systemd/juhradialmx-daemon.service
# JuhRadial MX Daemon - systemd user service
#
# This service manages the JuhRadial MX daemon which handles:
# - MX Master 4 gesture button detection via evdev
# - D-Bus communication with KWin overlay
# - Profile management and action execution
#
# Installation:
#   sudo install -Dm644 juhradialmx-daemon.service /usr/lib/systemd/user/juhradialmx-daemon.service
#   systemctl --user daemon-reload
#   systemctl --user reenable juhradialmx-daemon.service && systemctl --user start juhradialmx-daemon.service
#   (reenable, not enable: it drops stale [Install] symlinks from older versions)
#
# Management:
#   systemctl --user status juhradialmx-daemon   # Check status
#   systemctl --user restart juhradialmx-daemon  # Restart service
#   journalctl --user -u juhradialmx-daemon -f   # View logs

[Unit]
Description=JuhRadial MX Daemon - Radial menu for Logitech MX Master 4
Documentation=https://github.com/juhhally/juhradial-mx
# Ordering + lifecycle only. This unit must NEVER pull graphical-session.target
# in via Wants=/Requires=: that activated the target before the desktop's
# session manager registered, which GNOME treats as "a graphical session is
# already running" and GDM then kills the session - a hard login loop
# (issue #67). graphical-session.target is started by the session manager
# alone; we only order after it and stop with it. Being WantedBy targets
# (see [Install]) is safe: symlinks start US when a target activates, they
# never activate the target.
After=graphical-session.target
PartOf=graphical-session.target

# D-Bus session is required for communication with KWin
Requires=dbus.socket

# Limit restart attempts to prevent infinite restart loops.
# These are [Unit] keys; in [Service] systemd ignores them with a warning.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/juhradiald
# Restart on crashes (SIGSEGV, SIGABRT, etc.) and watchdog kills.
# on-abnormal = non-clean signals, watchdog, timeout — NOT SIGTERM/clean exit.
Restart=on-abnormal
RestartSec=5s

# Disable watchdog to prevent systemd drop-ins (e.g. 10-timeout-abort.conf
# on Fedora) from killing the daemon when it doesn't send sd_notify heartbeats.
WatchdogSec=0

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security hardening (user service, so limited options)
NoNewPrivileges=yes
# NOTE: PrivateTmp intentionally NOT set here. The daemon creates temporary
# KWin script files and passes their paths to KWin via D-Bus. With PrivateTmp
# the files would live in a private /tmp namespace invisible to KWin, causing
# the loadScript D-Bus call to silently fail (file not found on KWin's side).

# Resource limits
MemoryMax=100M
CPUQuota=50%

[Install]
# Both targets: graphical-session.target for systemd-managed desktops (GNOME,
# Plasma systemd boot, COSMIC, uwsm) gives correct ordering + lifecycle;
# default.target keeps autostart working on sessions that never activate
# graphical-session.target (bare exec-launched Hyprland/sway/niri/river, i3,
# XFCE, Cinnamon, Plasma with systemd boot disabled), where the old unit
# autostarted and graphical-session.target alone would silently never fire.
# Starting at default.target cannot recreate the #67 loop: without a Wants=
# line the unit does not activate graphical-session.target, it only starts
# earlier in the ordering-free case.
WantedBy=graphical-session.target default.target
