# packaging/systemd/juhradialmx-daemon.service
# JuhRadial MX Daemon - systemd user service
#
# This service manages the JuhRadial MX daemon which handles:
# - MX Master 4 gesture button detection via evdev
# - D-Bus communication with KWin overlay
# - Profile management and action execution
#
# Installation:
#   sudo install -Dm644 juhradialmx-daemon.service /usr/lib/systemd/user/juhradialmx-daemon.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now juhradialmx-daemon.service
#
# Management:
#   systemctl --user status juhradialmx-daemon   # Check status
#   systemctl --user restart juhradialmx-daemon  # Restart service
#   journalctl --user -u juhradialmx-daemon -f   # View logs

[Unit]
Description=JuhRadial MX Daemon - Radial menu for Logitech MX Master 4
Documentation=https://github.com/juhhally/juhradial-mx
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target

# D-Bus session is required for communication with KWin
Requires=dbus.socket

[Service]
Type=simple
ExecStart=/usr/local/bin/juhradiald
# Restart on crashes (SIGSEGV, SIGABRT, etc.) and watchdog kills.
# on-abnormal = non-clean signals, watchdog, timeout — NOT SIGTERM/clean exit.
Restart=on-abnormal
RestartSec=5s

# Limit restart attempts to prevent infinite restart loops
StartLimitIntervalSec=60
StartLimitBurst=5

# Disable watchdog to prevent systemd drop-ins (e.g. 10-timeout-abort.conf
# on Fedora) from killing the daemon when it doesn't send sd_notify heartbeats.
WatchdogSec=0

# Environment
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1

# Security hardening (user service, so limited options)
NoNewPrivileges=yes
# NOTE: PrivateTmp intentionally NOT set here. The daemon creates temporary
# KWin script files and passes their paths to KWin via D-Bus. With PrivateTmp
# the files would live in a private /tmp namespace invisible to KWin, causing
# the loadScript D-Bus call to silently fail (file not found on KWin's side).

# Resource limits
MemoryMax=100M
CPUQuota=50%

[Install]
WantedBy=default.target
